# src/etc/udev/rules.d/97-kobocloud.rules
KERNEL=="eth*", RUN+="/usr/local/kobocloud/udev_program.sh" 
KERNEL=="wlan*", RUN+="/usr/local/kobocloud/udev_program.sh" 
KERNEL=="lo", RUN+="/usr/local/kobocloud/udev_mount.sh"

